# systemd_units/sync_local.service
[Unit]
Description=Google Drive Sync: Local (inotify based) to Google Drive
After=network-online.target

[Service]
Type=simple
ExecStart=/full/path/to/gdrive-sync/sync_local_to_drive.sh
User=brandon
Group=brandon
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
