# systemd_units/sync_remote.service
[Unit]
Description=Google Drive Sync: Remote (Google Drive) to Local

[Service]
Type=oneshot
ExecStart=/full/path/to/gdrive-sync/sync_drive_to_local.sh
User=brandon
Group=brandon
